# autoarctest/rc.local
#!/bin/sh

mkdir -p /var/log/autoarctest
sleep 30
/home/jaadams/arctest/autoarctest/autotest.py > /var/log/autoarctest/$(date --iso-8601=seconds).log 2>&1
