# autoarctest/rc.local
#!/bin/sh

/home/jaadams/arctest/autoarctest/autotest.py > /var/log/autoarctest/$(date --iso-8601=seconds).log 2>&1
